# templates/oneshot-timer.timer.j2
[Unit]
Description=mpdPi project, run {{ unit }}

[Timer]
{{ timer }}

[Install]
WantedBy=timers.target
